# files/system/etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitInactiveSec=
RandomizedDelaySec=
RandomizedDelaySec=10m
OnCalendar=*-*-* 05:30:00 UTC
Persistent=true
