# etc/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnCalendar=Sat 20:00
Persistent=true
